# integration/openssh/modification/example/rsp-sshd.service
[Unit]
Description=RSP SSH Server Forwarder
Documentation=https://github.com/alienpopper/remote_socket_protocol
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust ExecStart path to match where rsp_sshd was installed.
ExecStart=/usr/local/bin/rsp_sshd /etc/rsp-sshd/rsp_sshd.conf.json
Restart=on-failure
RestartSec=5

# sshd -i requires root for privilege separation (default).
# For test environments with UsePrivilegeSeparation=no, a dedicated user suffices.
User=root

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rsp-sshd

TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
